# home/autostart/dot-config/autostart/steam.desktop
[Desktop Entry]
Actions=Store;Community;Library;Servers;Screenshots;News;Settings;BigPicture;Friends;
Categories=Network;FileTransfer;Game;
Comment[en_US]=Application for managing and playing games on Steam
Comment=Application for managing and playing games on Steam
Exec=/usr/bin/steam %U -silent
GenericName[en_US]=
GenericName=
Icon=steam
MimeType=
Name[en_US]=Steam
Name=Steam
Path=
PrefersNonDefaultGPU=true
StartupNotify=true
Terminal=false
TerminalOptions=
Type=Application
X-KDE-RunOnDiscreteGpu=true
X-KDE-SubstituteUID=false
X-KDE-Username=

[Desktop Action BigPicture]
Exec=/usr/bin/steam steam://open/bigpicture
Name=Big Picture

[Desktop Action Community]
Exec=/usr/bin/steam steam://url/CommunityHome/
Name=Community

[Desktop Action Friends]
Exec=/usr/bin/steam steam://open/friends
Name=Friends

[Desktop Action Library]
Exec=/usr/bin/steam steam://open/games
Name=Library

[Desktop Action News]
Exec=/usr/bin/steam steam://openurl/https://store.steampowered.com/news
Name=News

[Desktop Action Screenshots]
Exec=/usr/bin/steam steam://open/screenshots
Name=Screenshots

[Desktop Action Servers]
Exec=/usr/bin/steam steam://open/servers
Name=Servers

[Desktop Action Settings]
Exec=/usr/bin/steam steam://open/settings
Name=Settings

[Desktop Action Store]
Exec=/usr/bin/steam steam://store
Name=Store
